# stage3/03-leo-ros/files/systemd/ros-nodes.service
[Unit]
Description=ROS nodes
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/start-ros-nodes

[Install]
WantedBy=ros.target