# stage3/03-leo-ros/files/systemd/uros-agent.service
[Unit]
Description=Micro-ROS agent

[Service]
Type=simple
ExecStart=/usr/bin/start-uros-agent

[Install]
WantedBy=ros.target